# installer/linux/debian_package/etc/systemd/system/windscribe-helper.service
[Unit]
Description=Windscribe helper service
Before=network-pre.target
Wants=network-pre.target

[Service]
Type=simple
ExecStart=/usr/local/windscribe/helper

[Install]
WantedBy=multi-user.target