# installer/linux/debian_package/DEBIAN/postinst
#!/bin/bash
set -e
systemctl enable windscribe-helper
systemctl start windscribe-helper