# installer/linux/debian_package/DEBIAN/postinst
#!/bin/bash
# This file will be overwritten by rpm specific script which is passed to fpm in build script
# So if multi-platform change is needed remember to change that file as well
set -e
systemctl enable windscribe-helper
systemctl start windscribe-helper
echo linux_deb_x64 > ../etc/windscribe/platform # debian-specific

ln -s /usr/local/windscribe/windscribe-cli /usr/local/bin/windscribe-cli

