# installer/linux/debian_package/DEBIAN/prerm
#!/bin/bash
set -e
systemctl stop windscribe-helper
systemctl disable windscribe-helper

rm /usr/local/bin/windscribe-cli
