# src/installer/gui/linux/arch_package/windscribe.install
pre_install() {
        end=$((SECONDS+10))

        while :
        do
          if [ -z "$(ps -eo args | grep WindscribeEngine | grep -v grep)" ]; then
            break
          fi

          if [ $SECONDS -gt $end ]; then
            echo 'Error during Windscribe installation. WindscribeEngine is running. Please re-launch Windscribe. And try to update again.'
            exit 1
          fi
        done

        systemctl stop windscribe-helper > /dev/null 2>&1 || true
        systemctl disable windscribe-helper > /dev/null 2>&1 || true

        echo "Finish pre-install script"
}

post_install() {
        systemctl enable NetworkManager.service || true
        systemctl start NetworkManager.service || true
        systemctl enable systemd-resolved || true
        systemctl start systemd-resolved || true
        systemctl enable windscribe-helper
        systemctl start windscribe-helper
        ln -sf /opt/windscribe/windscribe-cli /usr/bin/windscribe-cli
        update-desktop-database || true
        setcap cap_setgid+ep /opt/windscribe/Windscribe
        mkdir -p /etc/windscribe
        echo linux_zst_x64 > /etc/windscribe/platform
}

pre_upgrade() {
        set -e
        systemctl stop windscribe-helper || true
        systemctl disable windscribe-helper || true
}

post_upgrade() {
        post_install
        rm -rf /etc/windscribe/stunnel.conf /etc/windscribe/config.ovpn /etc/windscribe/windscribe_servers
}

pre_remove() {
        killall -q Windscribe || true
        systemctl stop windscribe-helper || true
        systemctl disable windscribe-helper || true
        if [ -x /opt/windscribe/helper ]; then
            /opt/windscribe/helper --reset-mac-addresses || true
        fi
        userdel -f windscribe || true
        groupdel -f windscribe || true
        rm -rf /etc/windscribe
        rm -rf /var/log/windscribe
        rm -rf /var/lib/windscribe
        rm -f /usr/bin/windscribe-cli
}
