# src/installer/linux/cli/overlay/usr/lib/systemd/user/windscribe.service
[Unit]
Description=Windscribe service

[Service]
Type=simple
ExecStart=/opt/windscribe/Windscribe

[Install]
WantedBy=default.target
